# deploy/gunicorn.service
[Unit]
Description=Gunicorn instance to serve cyberlab-admin
After=network.target

[Service]
User=admin
Group=caddy
WorkingDirectory=/home/admin/Admin-Panel/cyberlab-admin
Environment="PATH=/home/admin/Admin-Panel/cyberlab-admin/.venv/bin"
Environment="FLASK_ENV=production"
EnvironmentFile=/home/admin/Admin-Panel/cyberlab-admin/.env
ExecStart=/home/admin/Admin-Panel/cyberlab-admin/.venv/bin/gunicorn --config gunicorn.conf.py wsgi:app
Restart=always
RestartSec=5
KillSignal=SIGQUIT
TimeoutStopSec=30
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
AmbientCapabilities=
CapabilityBoundingSet=
UMask=0007
RuntimeDirectory=cyberlab-admin
RuntimeDirectoryMode=0770
RuntimeDirectoryPreserve=no
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
